# ansible/roles/vps_service_foundation/templates/nutsnews-ops-portal-collector.service.j2
[Unit]
Description=NutsNews Operations Portal read-only data collector
Documentation={{ vps_service_foundation_docs_base_url }}
After=network-online.target docker.service
Wants=network-online.target
ConditionPathIsDirectory={{ vps_service_foundation_portal_data_dir }}

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=-{{ vps_service_foundation_free_tier_env_file }}
Environment="NUTSNEWS_ROOT_DIR={{ vps_service_foundation_root }}"
Environment="NUTSNEWS_PORTAL_OUTPUT={{ vps_service_foundation_portal_status_file }}"
Environment="NUTSNEWS_BACKUPS_DIR={{ vps_service_foundation_backups_dir }}"
Environment="NUTSNEWS_PUBLIC_IPV4={{ vps_service_foundation_public_ipv4 }}"
Environment="NUTSNEWS_PUBLIC_IPV6={{ vps_service_foundation_public_ipv6 }}"
Environment="NUTSNEWS_DOCS_BASE_URL={{ vps_service_foundation_docs_base_url }}"
Environment="NUTSNEWS_INFRA_REPO_URL={{ vps_service_foundation_infra_repo_url }}"
Environment="NUTSNEWS_DEPLOYED_COMMIT_FILE={{ vps_service_foundation_deployed_infra_commit_file }}"
Environment="NUTSNEWS_APPLY_STATUS_FILE={{ vps_service_foundation_apply_status_file }}"
Environment="NUTSNEWS_APP_ENABLED={{ vps_service_foundation_nutsnews_app_enabled | ternary('1', '0') }}"
Environment="NUTSNEWS_APP_ROUTE_ENABLED={{ vps_service_foundation_nutsnews_app_route_enabled | ternary('1', '0') }}"
Environment="NUTSNEWS_APP_CONTAINER_NAME={{ vps_service_foundation_nutsnews_app_container_name }}"
Environment="NUTSNEWS_APP_CONTAINER_PORT={{ vps_service_foundation_nutsnews_app_container_port }}"
Environment="NUTSNEWS_APP_IMAGE={{ vps_service_foundation_nutsnews_app_image }}"
Environment="NUTSNEWS_APP_IMAGE_REPO={{ vps_service_foundation_nutsnews_app_image_repo }}"
Environment="NUTSNEWS_APP_IMAGE_TAG={{ vps_service_foundation_nutsnews_app_image_tag }}"
Environment="NUTSNEWS_APP_HEALTH_PATH={{ vps_service_foundation_nutsnews_app_health_path }}"
Environment="NUTSNEWS_APP_ROUTE_PATH={{ vps_service_foundation_nutsnews_app_route_path }}"
Environment="NUTSNEWS_APP_APPLY_MARKER_FILE={{ vps_service_foundation_nutsnews_app_apply_marker_file }}"
Environment="NUTSNEWS_APP_ENV_FILE={{ vps_service_foundation_nutsnews_app_env_file }}"
Environment="NUTSNEWS_APP_SECRET_ENV_KEYS={{ vps_service_foundation_nutsnews_app_secret_env_keys | join(',') }}"
Environment="NUTSNEWS_APP_REQUIRED_SECRET_KEYS={{ vps_service_foundation_nutsnews_app_required_secrets | join(',') }}"
Environment="NUTSNEWS_REPORTING_STATUS_FILE={{ vps_service_foundation_reporting_status_file }}"
Environment="NUTSNEWS_BACKUP_STATUS_FILE={{ vps_service_foundation_backup_status_file }}"
Environment="NUTSNEWS_DISK_SCAN_CACHE_FILE={{ vps_service_foundation_disk_scan_cache_file }}"
Environment="NUTSNEWS_DISK_SCAN_CACHE_SECONDS={{ vps_service_foundation_disk_scan_cache_seconds }}"
Environment="NUTSNEWS_DISK_SCAN_ROOTS={{ vps_service_foundation_disk_scan_roots | join(',') }}"
Environment="NUTSNEWS_ALLOY_ENABLED={{ vps_service_foundation_grafana_alloy_enabled | ternary('1', '0') }}"
Environment="NUTSNEWS_ALLOY_COLLECT_DOCKER={{ vps_service_foundation_grafana_alloy_collect_docker | ternary('1', '0') }}"
Environment="NUTSNEWS_ALLOY_COLLECT_DOCKER_LOGS={{ vps_service_foundation_grafana_alloy_collect_docker_logs | ternary('1', '0') }}"
Environment="NUTSNEWS_ALLOY_SERVICE=alloy.service"
Environment="NUTSNEWS_ALLOY_READY_URL={{ vps_service_foundation_grafana_alloy_ready_url }}"
Environment="NUTSNEWS_ALLOY_ERROR_WINDOW={{ vps_service_foundation_grafana_alloy_error_validation_since }}"
Environment="NUTSNEWS_ALLOY_TEXTFILE_DIR={{ vps_service_foundation_grafana_alloy_textfile_dir }}"
Environment="HOME=/tmp"
Environment="DOCKER_CONFIG=/tmp/docker"
Environment="XDG_CACHE_HOME=/tmp/cache"
ExecStart={{ vps_service_foundation_collector_bin }}
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
ReadWritePaths={{ vps_service_foundation_portal_data_dir }} {{ vps_service_foundation_collector_log_dir }}

[Install]
WantedBy=multi-user.target
